# projects/dreadcam/root/vidmon.crontab
* * * * * /root/vidmon >> /root/info/vidmon.cron.log
